# debian/postinst
#!/bin/bash
# postinst script for flugbuch2

set -e

. /usr/share/debconf/confmodule
db_version 2.0 || [ $? - lt 30 ]


# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

#. /usr/share/debconf/confmodule
#. /usr/share/dbconfig-common/dpkg/postinst.pgsql
#dbc_go flightpred $@


case "$1" in
	configure)	
		# set up the environment for the webapp
		mkdir -p /usr/run/flugbuch2
    	chown -R www-data:www-data /usr/run/flugbuch2

		mkdir -p /usr/share/flugbuch2/data
		chown -R www-data:www-data /usr/share/flugbuch2/data

		mkdir -p /usr/share/flugbuch2/igc
		chown -R www-data:www-data /usr/share/flugbuch2/igc

		chown -R www-data:www-data /usr/share/flugbuch2/www_docroot


		# register with apache
		a2enmod fastcgi
		echo "register flugbuch2 with apache"
		a2ensite flugbuch2.local.conf
		a2ensite flugbuch2.ulrichard.ch.conf
	
		db_stop

		if which invoke-rc.d >/dev/null 2>&1; then
	     	invoke-rc.d apache2 restart || true
     	else
     		/etc/init.d/apache2 restart || true
     	fi
		echo "apache restarted"
		;;
	*)
		#
		# Unknown action - do nothing.
		#
		exit 0
		;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

